# proxy_auth_trigger.rules
SUBSYSTEM=="net", ACTION=="change", RUN+="/usr/bin/systemctl start tunet_proxy_auth.service"
